# supervisord.conf
[supervisord]
nodaemon=true

[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true

[program:telegram_bot]
command=python -m app.bot.main
autostart=true
autorestart=true
